# packaging/systemd/debian/postinst
#!/bin/bash

systemctl enable agent_protocol_splitter.service communication_link.service micrortps_agent.service px4_mavlink_ctrl.service

